# etc/supervisor/conf.d/authservice.ini
[program:auth_service_api]
command=/usr/bin/java -cp /opt/authservice/AuthService-1.0-SNAPSHOT-jar-with-dependencies.jar com.shareplaylearn.AuthService
