# etc/supervisor/conf.d/authservice.ini
[program:authservice]
user=authservice
command=/usr/bin/java -cp /opt/authservice/AuthService-1.1.0-SNAPSHOT-jar-with-dependencies.jar com.shareplaylearn.AuthService
