# etc/supervisor/conf.d/fileservice.ini
[program:fileservice]
user=fileservice
command=/usr/bin/java -Dlogback.configurationFile=/etc/fileservice/logback.xml -cp /opt/fileservice/FileService-1.1-SNAPSHOT-jar-with-dependencies.jar com.shareplaylearn.FileService
